# systemd/pl_boot.service
[Unit]
Description=pl_boot service

[Service]
ExecStart=/etc/init.d/pl_wrapper
Type=oneshot
