# systemd/pl_boot.service
[Unit]
Description=pl_boot service
Wants=network-online.target
After=network.target

# in debug mode we fork bash -x so that we can get a grip on
# what exactly is going on in here
[Service]
ExecStart=/bin/bash -x /etc/init.d/pl_sysinit > /dev/console
ExecStart=/bin/bash -x /etc/init.d/pl_boot > /dev/console
Type=oneshot
