# systemd/pl_boot.service
[Unit]
Description=pl_boot service

[Service]
ExecStart=/etc/init.d/pl_boot
ExecStartPre=/etc/init.d/pl_sysinit
Type=oneshot

